# xdg/config/systemd/user/edit-server.socket
[Unit]
Description=Edit Server

[Socket]
ListenStream=9292

[Install]
WantedBy=default.target

